# artifacts/raspberrypi/etc/cron.hourly/reconnect-wg-walbrix
#!/bin/sh
set -e
wb wg ping --success-if-not-active --quiet || systemctl restart wg-quick\@wg-walbrix
