# artifacts/nextcloud/etc/systemd/system/apache2.service.d/cloud-server-init.conf
[Unit]
After=mysqld.service

[Service]
ExecStartPre=/usr/sbin/cloud-server-init
