# artifacts/nextcloud/etc/systemd/system/apache2.service.d/nextcloud-init.conf
[Unit]
After=mysqld.service

[Service]
ExecStartPre=/usr/bin/nextcloud-init
