# artifacts/whitebase/lib/systemd/system/generate-wg-walbrix-key.service
[Unit]
Description=Generate wg-walbrix secret key
Requires=local-fs.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/wb wg genkey

[Install]
WantedBy=multi-user.target
