# artifacts/owncloud/etc/cron.hourly/owncloud-cron
#!/bin/sh
set -e

LOGFILE=/var/log/owncloud-cron.log
tty -s || exec >$LOGFILE 2>&1

occ system:cron
