# artifacts/walbrix/etc/systemd/system/kmsconvt@tty1.service.d/walbrix.conf
[Service]
ExecStart=
ExecStart=/usr/bin/kmscon --login "--vt=%I" --seats=seat0 --no-switchvt -- /usr/sbin/wb login
Restart=always
