# packages/@cloudflared/usr/lib/systemd/system/cloudflared.service
[Unit]
Description=cloudflared
After=network.target

[Service]
TimeoutStartSec=0
Type=notify
ExecStart=cloudflared --no-autoupdate --config /etc/cloudflared/config.yml tunnel --edge-ip-version auto run
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
