# artifacts/owncloud/etc/systemd/system/apache2.service.d/owncloud-init.conf
[Unit]
After=mysqld.service

[Service]
ExecStartPre=/usr/bin/owncloud-init
TimeoutStartSec=infinity
