# artifacts/raspberrypi/etc/systemd/system/wg-quick@wg-walbrix.service.d/override.conf
[Unit]
Requires=systemd-resolved.service
StartLimitIntervalSec=10
[Service]
Restart=on-failure
StartLimitBurst=10
