# artifacts/whitebase/usr/lib/systemd/system/installer.target
[Unit]
Description=Installer Mode
Requires=graphical.target
