# packages/dev-db/mysql/etc/systemd/system/mysqld.service.d/copy-up-to-overlay.conf
[Unit]
Requires=mysql-copy-up-to-overlay.service

