# packages/dev-db/mysql/etc/systemd/system/mysql-copy-up-to-overlay.service
[Unit]
Description=Copy-up MySQL's data directory to overlay
Before=mysqld.service

[Service]
ExecStart=/usr/sbin/mysql-copy-up-to-overlay
Type=oneshot
RemainAfterExit=no
